# etc/services.d/nginx/run
#!/bin/sh

nginx -g "daemon off;" -c /dev/shm/nginx.conf
